# project/cookbook/mongodb/templates/default/supervisor.ini.erb
[program:mongodb]
command=bash -c "mongod -f <%=@config_dir%>/mongod.conf"
process_name=%(program_name)s